# crontabs/purge_cached_files
0-59/20 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p purgefiles $ENSTORE_DIR/purge_files_wrapper >/dev/null 2>/dev/null
